# CICD-configs/Docker/supervisord.conf
[supervisord]
nodaemon=true
[program:sonarqube]
command:/opt/sonarqube/bin/run.sh
[program:devnill]
command:ENTRYPOINT tail -F /dev/null
[program:sshd]
command:/etc/init.d/sshd start
